# etc/supervisor.d/dnsmasq.ini
[program:dnsmasq]
command=/usr/sbin/dnsmasq --no-daemon
autorestart=true