# etc/supervisor.d/hostapd.ini
[program:hostapd]
command=/usr/sbin/hostapd -P /var/run/hostapd.pid /etc/hostapd.conf
autorestart=true
