# etc/supervisor.d/igmpproxy.ini
[program:igmpproxy]
command=/usr/sbin/igmpproxy -n /etc/igmpproxy.conf
autorestart=true