# etc/supervisor.d/radvd.ini
[program:radvd]
command=/usr/sbin/radvd -p /var/run/radvd.pid -n
autorestart=true
